# tile-processor/crontab
SHELL=/bin/bash

* * * * * root date > /data/tiles/health-check 2>&1 | logger -t "health-check"
0 0 * * * root /app/get-and-process-tiles.sh | logger -t "get-and-process-tiles"
